# systemd/dashboard-gunicorn.socket
[Unit]
Description=dashboard gunicorn socket

[Socket]
ListenStream=/run/dashboard-gunicorn/socket
User=www-data

[Install]
WantedBy=sockets.target
